# sv/cloud9/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

exec 2>&1
exec node /core/server.js -w /docker
